# certbot_renew/crontab
# Перевыпуск истекших SSL-сертификатов каждый день в 23:00
0 23 * * *   root    certbot renew --quiet
